# Inventory Submissions/systemd/inventory-automation.service
[Unit]
Description=Inventory Automation daily run
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/workspaces/Inventory-Automation
ExecStart=/workspaces/Inventory-Automation/run_daily.sh

[Install]
WantedBy=multi-user.target